# deploy/celery/celery-ciptea-ia.service
[Unit]
Description=Celery IA Worker Service for CIPTEA
After=network.target redis-server.service

[Service]
Type=simple
User=gabinete
Group=www-data
WorkingDirectory=/var/www/ciptea
EnvironmentFile=/var/www/ciptea/.env
Environment=DJANGO_SETTINGS_MODULE=core.settings
ExecStart=/var/www/ciptea/venv/bin/celery -A core worker -l info -Q ia_tasks --concurrency=2 --max-tasks-per-child=10 -n worker_ia@%%h
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
